# deploy/systemd/bottwitch-web.service
[Unit]
Description=Panel web Bottwitch
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/bottwitch
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/bottwitch/scripts/linux/start_web_panel.sh
Restart=always
RestartSec=5
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
